# build/src/izpack/edge-server.conf
description "This is an upstart job file for RSNA ISN Edge Server UI"

start on started networking
stop on runlevel [!2345]
respawn

limit nofile 4096 4096

script
bash << "EOF"
  su - edge
  . /etc/rsna.conf
  export RSNA_ROOT
  export SECRET_KEY_BASE
  export OPENAM_URL
  %{rsna.root}/torquebox-%{torquebox.version}/jboss/bin/standalone.sh -b %{server.host} -Dhttp.port=3000 > /dev/null 2>&1
EOF
end script
